# root/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

/usr/sbin/crond -c /etc/crontabs